# systemd/hoera10jaar-update-urls.timer
[Unit]
Description=update spaceapi urls from directory

[Timer]
OnCalendar=daily

[Install]
WantedBy=hoera10jaar.target
