# systemd/hoera10jaar-update-open.timer
[Unit]
Description=update space states for hoera10jaar

[Timer]
OnCalendar=minutely

[Install]
WantedBy=hoera10jaar.target
